# systemd/backup-drives.timer
[Unit]
Description=Run Achterhus Backup Daily

[Timer]
OnCalendar=*-*-* 03:11:00
Persistent=true

[Install]
WantedBy=default.target
